# supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=debug                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:harness]
command=/harness/bootstrap 
directory=/harness
autostart=true

[program:crs]
command=/harness/crs/crs/app.py -p 56789 
directory=/harness/crs/crs
autostart=false 

[program:irm-net]
command=/harness/irm-net/irm-net/app.py -p 7779 -o localhost -t 56789
directory=/harness/irm-net
autostart=false

[program:irm-nova]
command=/harness/irm-nova/irm-nova.py -c ../cfg/irm-nova.cfg 
directory=/harness/irm-nova
autostart=false

[program:hresmon]
command=/harness/irm-nova/hresmon.py
directory=/harness/irm-nova
autostart=false

[program:irm-neutron]
command=/harness/irm-neutron/irm-neutron.py -c ../cfg/irm-neutron.cfg
directory=/harness/irm-neutron
autostart=false
